# target_aws_common/usr/lib/systemd/system/dcos-01-startup.service
[Unit]
Description=Creates DCOS environment in Docker Container

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/start.sh

[Install]
WantedBy=default.target
